# RPi/serialread.conf
[program:serialread]
command=/home/pi/bin/serialread.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/serialread.log
stdout_logfile=/var/log/serialread.log
